# runners/dev.tquick.act-runner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- macOS LaunchAgent for act_runner (MacBook)
     Install: cp dev.tquick.act-runner.plist ~/Library/LaunchAgents/
              launchctl load ~/Library/LaunchAgents/dev.tquick.act-runner.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.tquick.act-runner</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/act_runner</string>
        <string>daemon</string>
        <string>--config</string>
        <string>/Users/tquick/.config/act_runner/config.yaml</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/act_runner.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/act_runner.stderr.log</string>
</dict>
</plist>
